# g3x-processor@.service
# G3X SD Card Processor Service
# /lib/systemd/system/g3x-processor@.service
#
# This is a template service that gets instantiated with the device name
# when an SD card is inserted.

[Unit]
Description=G3X Flight Data Processor for %i
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
Type=oneshot
RemainAfterExit=no

# Mount the device, process, then unmount
ExecStartPre=/bin/mkdir -p /mnt/g3x_sdcard
ExecStartPre=/bin/mount -o ro,sync /dev/%i /mnt/g3x_sdcard

ExecStart=/usr/local/bin/g3x_processor.py /mnt/g3x_sdcard --config /etc/g3x_processor/config.json

ExecStopPost=/bin/umount /mnt/g3x_sdcard
ExecStopPost=/bin/rmdir /mnt/g3x_sdcard

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=g3x-processor

# Safety: timeout after 10 minutes
TimeoutStartSec=600

[Install]
WantedBy=multi-user.target
