# systemd/avcardtool-processor@.service
[Unit]
Description=AVCardTool Processor for %i
Documentation=https://github.com/elvinzhou/g3xuploader
BindsTo=dev-%i.device
After=dev-%i.device

[Service]
Type=oneshot
ExecStart=/usr/local/bin/avcardtool auto-process /dev/%i
StandardOutput=journal
StandardError=journal
User=AVCARDTOOL_USER
Group=AVCARDTOOL_USER

# Security hardening
PrivateTmp=yes
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=AVCARDTOOL_DATA_DIR AVCARDTOOL_CONFIG_DIR

# Resource limits
TimeoutStartSec=600
CPUQuota=50%
MemoryMax=512M

[Install]
WantedBy=multi-user.target
